# x/systemd/user/echo-pwd.service
[Unit]
Description=echo pwd

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'pwd'

[Install]
WantedBy=
