# x/systemd/user/feh-wallpaper.timer
[Unit]
Description=rotate feh-wallpaper.service

[Timer]
OnUnitActiveSec=2min
Unit=feh-wallpaper.service

[Install]
WantedBy=
