# x/systemd/user/feh-wallpaper.timer
[Unit]
Description=for feh-wallpaper.service

[Timer]
OnUnitActiveSec=1m
Unit=feh-wallpaper.service

[Install]
WantedBy=graphical-session.target
