# x/systemd/user/getter.timer
[Unit]
Description=wall getter timer

[Timer]
OnUnitActiveSec=2h 20min
Unit=getter.service

[Install]
WantedBy=
