# x/systemd/user/sway-getter.timer
[Unit]
Description=wall getter timer for sway

[Timer]
OnUnitActiveSec=2h 20min
Unit=sway-getter.service

[Install]
WantedBy=
