# x/systemd/user/sway-wallpaper.timer
[Unit]
Description=for sway-wallpaper.service

[Timer]
OnUnitActiveSec=1h
Unit=sway-wallpaper.service

[Install]
WantedBy=
